# volumio-controller-service/volumio-controller.service
[Unit]
Description = Volumio Controller

[Service]
ExecStart=/usr/bin/node /data/INTERNAL/programs/volumio-controller/index.js
Restart=always
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=volumio-controller
User=volumio
Group=volumio
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target
