# volumio-kiosk-service/volumio-kiosk.service
[Unit]
Description=Volumio Kiosk (static web UI for Volumio)
After=network-online.target
Wants=network-online.target

[Service]
ExecStart=/usr/bin/python3 -m http.server 8080 --directory /data/INTERNAL/programs/volumio-kiosk
Restart=on-failure
RestartSec=3
User=volumio
Group=volumio

[Install]
WantedBy=multi-user.target
